# infra/yt-backend.service
[Unit]
Description=YT Editor FastAPI Backend
After=network.target

[Service]
Type=simple
User=yt-editor
Group=yt-editor
WorkingDirectory=/opt/yt-editor/backend
EnvironmentFile=/opt/yt-editor/backend/.env
ExecStart=/opt/yt-editor/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
Environment=PATH=/opt/yt-editor/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/yt-editor/.cache/ms-playwright

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/yt-editor/data /opt/yt-editor/backend/config
PrivateTmp=true

[Install]
WantedBy=multi-user.target
